# CRAN/QA/BDR/scripts/crontab
5 6 * * * /data/gannet/ripley/R/packages/check-gcc > /data/gannet/ripley/R/packages/check-gcc.log 2>&1
48 7 * * * /data/gannet/ripley/R/packages/Rlog_diffs2
