# CRAN/QA/BDR/Solaris/x86/crontab
5 4,7,10,16,22 * * * /home/ripley/R/packages/check-swift > /home/ripley/R/packages/check-swift.log 2>&1
5 1 * * * /home/ripley/R/packages/tl-update > /home/ripley/R/packages/tl-update.log  2>&1
